# com.user.remindersync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.remindersync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/rabu/Projects/sync reminders/run_sync.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>StandardOutPath</key>
    <string>/Users/rabu/Projects/sync reminders/sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/rabu/Projects/sync reminders/sync.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/rabu/Projects/sync reminders</string>
</dict>
</plist>
